# templates/bamboo.service.erb
[Unit]
Description=Atlassian Bamboo Server Service
After=syslog.target network.target

[Service]
Type=forking
User=<%= scope['bamboo::user'] %>
ExecStart=<%= scope['bamboo::webappdir'] %>/bin/start-bamboo.sh
ExecStop=<%= scope['bamboo::webappdir'] %>/bin/stop-bamboo.sh

[Install]
WantedBy=multi-user.target
